# root/etc/services.d/server/run
#!/usr/bin/with-contenv bash

cd /app/synclounge || exit

exec \
	s6-setuidgid abc npm run server
